# devops/vt.service
[Unit]
Description=Vetements service

[Service]
ExecStart=~/projets/script/vt.sh

[Install]
WantedBy=multi-user.target